# ops/launchd/com.atlas.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.atlas.backup</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/jayakrishnakaryampudi/Documents/atlas/ops/backup.sh</string>
  </array>
  <!-- 10:30 local — after the daily pipeline, so the dump includes today's run -->
  <key>StartCalendarInterval</key>
  <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>30</integer></dict>
  <key>StandardOutPath</key><string>/Users/jayakrishnakaryampudi/Library/Logs/atlas-backup.log</string>
  <key>StandardErrorPath</key><string>/Users/jayakrishnakaryampudi/Library/Logs/atlas-backup.log</string>
</dict>
</plist>
